# launchd/com.danidin.pokemon-auctions.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.danidin.pokemon-auctions.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/danidin/Desktop/pokemon/src/run-daily.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/danidin/Desktop/pokemon</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/danidin/Desktop/pokemon/logs/daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/danidin/Desktop/pokemon/logs/daily.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
