# launchd/com.danidin.pokemon-auctions.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.danidin.pokemon-auctions.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/danidin/Desktop/pokemon/src/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/danidin/Desktop/pokemon</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/danidin/Desktop/pokemon/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/danidin/Desktop/pokemon/logs/server.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
